# app/deploy/systemd/writ-server.service
[Unit]
Description=Writ API (Express)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/app/server
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=PATH=/home/ubuntu/.nvm/versions/node/v24.13.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/etc/app.env
ExecStart=/home/ubuntu/.nvm/versions/node/v24.13.0/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
